# roles/openhab/templates/openhab.service
[Service]
User=
User=pi
Group=
Group=pi
